# deb/install/prerm
#!/bin/bash

systemctl stop waggle-upload-sync.timer || true
systemctl disable waggle-upload-sync.timer || true
systemctl stop waggle-upload-sync.service || true
